# deploy/paintsource.service
[Unit]
Description=PaintSource PRO (FastAPI)
After=network.target

[Service]
User=root
WorkingDirectory=/opt/paintsource_pro
ExecStart=/opt/paintsource_pro/.venv/bin/python app.py --host 127.0.0.1 --port 7100 --cols 80 --rows 24 --scale 10 --data /opt/paintsource_pro/data --autosave-sec 3
Restart=always
RestartSec=2
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
